# s6-overlay/s6-rc.d/langchain-rag/finish
#!/bin/sh
# LangChain RAG Service Finish Script

echo "[s6-langchain-rag] Shutting down LangChain RAG service..."
