# s6-overlay/s6-rc.d/langchain-rag/run
#!/bin/sh
# LangChain RAG Service v2.0 - Full Implementation
# Port: 3032
# Features: Document Ingestion, Smart Sync, READ_WRITE Mode

echo "[s6-langchain-rag] Starting LangChain RAG Service v2.0..."

cd /app/mini-services/langchain-rag-service

export PORT="${LANGCHAIN_RAG_PORT:-3032}"
export MEDICAL_RAG_URL="${MEDICAL_RAG_URL:-http://localhost:3031}"
export ZAI_API_KEY="${ZAI_API_KEY:-}"
export LOG_LEVEL="${LOG_LEVEL:-INFO}"

# Activate virtual environment if available
if [ -d "venv" ]; then
    echo "[s6-langchain-rag] Using virtual environment"
    . venv/bin/activate
fi

echo "[s6-langchain-rag] Starting on port ${PORT} (READ_WRITE mode)..."
echo "[s6-langchain-rag] Medical RAG URL: ${MEDICAL_RAG_URL}"
exec python index_full.py
