# s6-overlay/s6-rc.d/medical-rag/finish
#!/bin/sh
# Medical RAG Service Finish Script

echo "[s6-medical-rag] Shutting down Medical RAG service..."
