# s6-overlay/s6-rc.d/medical-rag/run
#!/bin/sh
# Medical RAG Service v2.0 - Full Implementation
# Port: 3031
# Features: Z.ai LLM, Semantic Search, Medical Knowledge Base

echo "[s6-medical-rag] Starting Medical RAG Service v2.0..."

cd /app/mini-services/medical-rag-service

export PORT="${MEDICAL_RAG_PORT:-3031}"
export ZAI_API_KEY="${ZAI_API_KEY:-}"
export LOG_LEVEL="${LOG_LEVEL:-INFO}"

# Activate virtual environment if available
if [ -d "venv" ]; then
    echo "[s6-medical-rag] Using virtual environment"
    . venv/bin/activate
fi

echo "[s6-medical-rag] Starting on port ${PORT}..."
exec python index_full.py
