# s6-overlay/s6-rc.d/caddy/finish
#!/bin/sh
# Caddy Service Finish Script

echo "[s6-caddy] Shutting down Caddy reverse proxy..."
caddy stop 2>/dev/null || true
