# s6-overlay/s6-rc.d/medasr/finish
#!/bin/sh
# MedASR Service Finish Script

echo "[s6-medasr] Shutting down MedASR service..."
