# s6-overlay/s6-rc.d/medasr/run
#!/bin/sh
# MedASR Service v3.0 - Medical Speech Recognition
# Port: 3033
# Features: VAD, Noise Reduction, Phonetic Matching, Negation Detection
# Dependencies: init-first only (standalone service)

echo "[s6-medasr] Starting MedASR v3.0 speech recognition service..."

cd /app/mini-services/medasr-service

# Set environment
export PORT="${MEDASR_PORT:-3033}"
export HF_TOKEN="${HF_TOKEN:-}"
export MODEL_CACHE_DIR="${MODEL_CACHE_DIR:-/app/models}"

# Check if dependencies are ready
if [ -f "/run/s6/basedir/env/SERVICE_READY" ]; then
    echo "[s6-medasr] ✅ Init complete, starting service"
fi

# Activate virtual environment if available (for VAD support)
if [ -d "venv" ]; then
    echo "[s6-medasr] Using virtual environment with VAD support"
    . venv/bin/activate
else
    echo "[s6-medasr] Warning: venv not found, VAD may not be available"
fi

echo "[s6-medasr] Starting uvicorn on port ${PORT}..."
exec uvicorn index:app --host 0.0.0.0 --port ${PORT} --workers 1
