# s6-overlay/s6-rc.d/db-backup/finish
#!/bin/sh
# Database Backup Service Finish Script

echo "[db-backup] Shutting down backup service..."
