# s6-overlay/s6-rc.d/nextjs/finish
#!/bin/sh
# Next.js Service Finish Script

echo "[s6-nextjs] Received shutdown signal..."
echo "[s6-nextjs] Graceful shutdown complete"
